# hash256-master.service
[Unit]
Description=HASH256 master miner
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
WorkingDirectory=/opt/hash256-miner
EnvironmentFile=/opt/hash256-miner/.env.master
ExecStart=/usr/bin/node /opt/hash256-miner/hash256-mine.mjs master
Restart=always
RestartSec=5
User=root

[Install]
WantedBy=multi-user.target
